# trogdor/auv-docker.service
[Unit]
Description=Start Docker container

[Service]
User=software
ExecStart=/home/software/docker-helper.sh run

[Install]
WantedBy=docker.target
