# install/runit/trogdor/run
#!/bin/sh
exec /home/software/cuauv/workspaces/worktrees/master/trogdor/trogdor_boot.sh